# DockerConfig/Queue/config/schedule.conf
[program:schedule]
command=/usr/local/bin/schedule.sh
autostart=true
autorestart=true
#stdout_logfile=/var/log/supervisord/schedule-stdout.log
#stderr_logfile=/var/log/supervisord/schedule-stderr.log